# pi-voice-assistant/backend/voice-assistant.service
[Unit]
Description=BMO/Sam Voice Assistant Backend
After=network.target

[Service]
User=sam
WorkingDirectory=/home/sam/oogway/pi-voice-assistant/backend
ExecStart=/home/sam/oogway/pi-voice-assistant/backend/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=inherit
StandardError=inherit

[Install]
WantedBy=multi-user.target
